# thinkpad-silent-fans/silent.service
[Unit]
Description=Silent Fan Control

[Service]
ExecStart=/usr/local/libexec/silent/init.sh

[Install]
WantedBy=multi-user.target
